# crontab
*/15 * * * * /home/dimo/dimo-control/do_it.sh